# lyzwlkjvip/mnt/pokemon/deploy_dev/supervisord.dir/bk/3_game_server.ini
[program:3_game_server]
directory=../release/
command=/usr/bin/python game_server.py game.cn.3
autostart=false
autorestart=false
priority=109
redirect_stderr=true
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=999999
