# lyzwlkjvip/mnt/pokemon/deploy_dev/supervisord.dir/bk/3_pvp_server.ini
[program:3_pvp_server]
directory=../release/pvp
command=%(here)s/../release/bin/host_server -config=pvp.cn.3
autostart=false
autorestart=false
priority=108
redirect_stderr=true
