# lyzwlkjvip/mnt/pokemon/deploy_dev/supervisord.dir/bk/3_storage_server.ini
[program:3_storage_server]
directory=../release/storage
command=%(here)s/../release/bin/host_server -config=storage.cn.3
autostart=false
autorestart=false
priority=107
redirect_stderr=true
