# lyzwlkjvip/mnt/pokemon/deploy_dev/supervisord.dir/bk/4_pvp_server.ini
[program:4_pvp_server]
directory=../release/pvp
command=%(here)s/../release/bin/host_server -config=pvp.cn.4
autostart=false
autorestart=false
priority=102
redirect_stderr=true