# lyzwlkjvip/mnt/pokemon/deploy_dev/supervisord.dir/account_db_server.ini
[program:account_db_server]
directory=../release/storage
command=%(here)s/../release/bin/host_server -config=account.cn.1
autostart=true  
autorestart=false
priority=101
redirect_stderr=true

