# lyzwlkjvip/mnt/pokemon/deploy_dev/supervisord.dir/anti_cheat_server.ini
[program:anti_cheat_server]  
directory=../release/anti_cheat
command=%(here)s/../release/bin/anti_cheat_server -config=agent.cn.1 -limit=1
autostart=true  
autorestart=false
priority=101
redirect_stderr=true

