# lyzwlkjvip/mnt/pokemon/deploy_dev/supervisord.dir/autochess_room_server.ini
[program:autochess_room_server]
directory=../release/autochess_room
command=%(here)s/../release/bin/autochess_room_server -config=autochess_room.cn.1
autostart=true  
autorestart=false
priority=102
redirect_stderr=true
