# lyzwlkjvip/mnt/pokemon/deploy_dev/supervisord.dir/comment_server.ini
[program:comment_server]
directory=../release/cross
command=%(here)s/../release/bin/host_server -config=comment.cn.1
autostart=true  
autorestart=false
priority=101
redirect_stderr=true