# lyzwlkjvip/mnt/pokemon/deploy_dev/supervisord.dir/crosslobby_server.ini
[program:crosslobby_server]  
directory=../release/crosslobby
command=%(here)s/../release/bin/host_server -config=crosslobby.cn.1
autostart=true  
autorestart=false
priority=102
redirect_stderr=true

