# lyzwlkjvip/mnt/pokemon/deploy_dev/supervisord.dir/cross_server.ini
[program:cross_server]  
directory=../release/cross
command=%(here)s/../release/bin/host_server -config=cross.cn.1 
autostart=true  
autorestart=false
priority=102
redirect_stderr=true

