# lyzwlkjvip/mnt/pokemon/deploy_dev/supervisord.dir/gm_server.ini
[program:gm_server]  
directory=../release/  
command=/usr/bin/python gm_server.py gm.cn.1
autostart=true  
autorestart=false
priority=105
redirect_stderr=true

