# lyzwlkjvip/mnt/pokemon/deploy_dev/supervisord.dir/login_server_1.ini
[program:login_server1]
directory=%(here)s/../release/login
command=%(here)s/../release/bin/login_server -config=login.cn.1
autostart=true
autorestart=false
priority=110
redirect_stderr=true
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=999999
