# lyzwlkjvip/mnt/pokemon/deploy_dev/supervisord.dir/mongodb_server.ini
[program:mongod]
command=mongod --config /www/server/mongodb/config.conf
autostart=false
autorestart=false
stderr_logfile=../deploy_dev/childlog/mongodb/mongod.err.log
stdout_logfile=../deploy_dev/childlog/mongodb/mongod.out.log
