# lyzwlkjvip/mnt/pokemon/deploy_dev/supervisord.dir/nsqlookupd.ini
[program:nsqlookupd]  
directory=./nsq
command=/mnt/pokemon/deploy_dev/nsq/nsqlookupd --http-address=127.0.0.1:4161 --tcp-address=127.0.0.1:4160
autostart=true
autorestart=true
priority=10
redirect_stderr=true

