# lyzwlkjvip/mnt/pokemon/deploy_dev/supervisord.dir/online_fight_forward_server.ini
[program:online_fight_forward_server]
directory=../release/online_fight_forward
command=%(here)s/../release/bin/online_fight_forward_server -config=online_fight_forward.cn.1
autostart=true  
autorestart=false
priority=101
redirect_stderr=true