# lyzwlkjvip/mnt/pokemon/deploy_dev/supervisord.dir/pvp_server_2.ini
[program:pvp_server2]
directory=../release/pvp
command=%(here)s/../release/bin/host_server -config=pvp.cn.2
autostart=true
autorestart=false
priority=102
redirect_stderr=true
