# lyzwlkjvip/mnt/pokemon/deploy_dev/supervisord.dir/pvp_server_1.ini
[program:pvp_server1]
directory=../release/pvp
command=%(here)s/../release/bin/host_server -config=pvp.cn.1
autostart=true
autorestart=false
priority=102
redirect_stderr=true
