# lyzwlkjvip/mnt/pokemon/deploy_dev/supervisord.dir/storage_server_2.ini
[program:storage_server2]
directory=../release/storage
command=%(here)s/../release/bin/host_server -config=storage.cn.2
autostart=true
autorestart=false
priority=101
redirect_stderr=true
