# lyzwlkjvip/mnt/pokemon/deploy_dev/supervisord.dir/transfer_daemon.ini
[program:transfer_daemon]
directory=../release/
command=/usr/bin/python transfer_server.py transfer_daemon
autostart=true
autorestart=true
priority=200
redirect_stderr=true
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=999999

