# lyzwlkjvip/mnt/pokemon/deploy_dev/childlog/bk/mongodb_server.ini
[program:mongod]
command=mongod --config /www/server/mongodb/config.conf
autostart=true
autorestart=true
stderr_logfile=/mnt/pokemon/deploy_dev/childlog/mongodb/mongod.err.log
stdout_logfile=/mnt/pokemon/deploy_dev/childlog/mongodb/mongod.out.log
